# emacs-integration.service
[Unit]
Description=Emacs remote integration service

[Service]
Environment=DISPLAY=:0
ExecStart=%h/go/bin/emacs-integration
