# emacs-integration.socket
[Unit]
Description=Emacs remote integration service activation socket

[Socket]
ListenStream=127.0.0.1:9999

[Install]
WantedBy=sockets.target
